# services/brain-autostart/brain-monitor.service
[Unit]
Description=BFMC Brain Monitor
After=network.target
Wants=network-online.target
ConditionPathExists=/opt/brain-autostart/monitor-dashboard.sh

[Service]
Type=simple
User=ewolf
Group=ewolf
WorkingDirectory=/home/ewolf/Documents/Brain
Environment=HOME=/home/ewolf
EnvironmentFile=/opt/brain-autostart/config.env
ExecStart=/opt/brain-autostart/monitor-dashboard.sh
Restart=always
RestartSec=10
StandardOutput=append:/var/log/brain-monitor.log
StandardError=append:/var/log/brain-monitor.log

[Install]
WantedBy=multi-user.target
